# homelab/ansible/files/net-watchdog/net-isolation-watchdog.service
[Unit]
Description=Network-isolation watchdog (self-heal a wedged-NIC node; safe vs gateway/switch outage)
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/net-isolation-watchdog.sh
