# homelab/ansible/files/net-watchdog/net-isolation-watchdog.timer
[Unit]
Description=Run the network-isolation watchdog every 30s
[Timer]
OnBootSec=90
OnUnitActiveSec=30
AccuracySec=5
[Install]
WantedBy=timers.target
